# bin/kmkr-arp-check.cron
#
# cron-jobs for munin
#

MAILTO=<<<mailto address>>>

*/5 * * * *     root test -x /opt/kmkr-arp-check/bin/kmkr-arp-check && /opt/kmkr-arp-check/bin/kmkr-arp-check "<<<mailto address>>>"
